# backend/infra/systemd/rightbricks-backend.service
[Unit]
Description=RightBricks Backend (Next.js)
After=network.target docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/rightbricks/backend
ExecStart=/usr/bin/docker compose -f /opt/rightbricks/backend/infra/docker/docker-compose.yml up -d --build
ExecStop=/usr/bin/docker compose -f /opt/rightbricks/backend/infra/docker/docker-compose.yml down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
